# ops/systemd/nanoclaw.service
[Unit]
Description=NanoClaw service (user)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/clawrocket
ExecStart=%h/projects/clawrocket/ops/systemd/start-nanoclaw.sh
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30
Environment=NODE_ENV=production
Environment=CLAWROCKET_SELF_RESTART=1
# Optional: uncomment and set if you keep runtime env in a separate file.
# EnvironmentFile=%h/projects/clawrocket/.env

[Install]
WantedBy=default.target
